# deploy/systemd/guiguzi-gateway.service
[Unit]
Description=Guiguzi Multi-Channel Gateway
Documentation=https://github.com/guiguzi/guiguzi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=guiguzi
Group=guiguzi
ExecStart=/usr/bin/node /opt/guiguzi/packages/nova-gateway/dist/index.js
WorkingDirectory=/opt/guiguzi
Restart=on-failure
RestartSec=5

# Environment
EnvironmentFile=/etc/guiguzi/env
Environment=NODE_ENV=production

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/guiguzi /var/log/guiguzi

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
